# supervisord.conf
[program:mysql]
command = /usr/bin/mysqld_safe > /dev/null 2>&1 &
autorestart = true

[program:apache]
command = /usr/sbin/httpd -k start
autorestart = true